# sharedemos-dev-two/supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf
user=postgres

[program:sshd]
command=/usr/sbin/sshd -D

[program:rabbitmq]
command=service rabbitmq-server start

[program:redis-server]
command=redis-server --daemonize yes